# systemd-units/update-unbound-ads.service
[Unit]
Description=Update Unbound ad blocking list (OISD Full — 127.0.0.1)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update-unbound-ads.sh

[Install]
WantedBy=multi-user.target
